# opt/implant/services/power-monitor.service
[Unit]
Description=Power Monitor
Wants=power-monitor.timer

[Service]
Type=oneshot
ExecStart=/opt/implant/scripts/power-monitor.sh
User=root
